# lb-haproxy-study/LB-HAProxy-service/haproxy.reload.path
[Unit]
Description=Monitor the file stopping file for changes, touch is enough

[Path]
PathChanged=/haproxy/haproxy.cfg
Unit=haproxy.reload.service

[Install]
WantedBy=multi-user.target
